# 99-custom-usb.rules
ACTION=="add", SUBSYSTEM=="block", DRIVER=="", RUN+="/git/devnotify/new.sh $attr{size}"
