# 99-imager-notify.rules
ACTION=="add|change", SUBSYSTEM=="block", DRIVER=="", RUN+="/root/imager/notify.sh $attr{size}"
