# 99-imager-notify.rules
ACTION=="add|change|remove", SUBSYSTEM=="block", DRIVER=="", RUN+="/root/notify http://localhost:8080 $attr{size}"
